# nuclear-camera.service
[Unit]
Description=Nuclear Camera — RTSP/HTTP Vision Ingestion
After=network.target fastvlm.service

[Service]
User=crew
WorkingDirectory=/home/crew/git/nuclear-camera
ExecStart=/home/crew/.venv/nuclear/bin/uvicorn main:app --host 0.0.0.0 --port 8092
Restart=on-failure
RestartSec=10s
Environment=CAMERA_CONFIG=/home/crew/git/nuclear-camera/cameras.toml

[Install]
WantedBy=multi-user.target
