# deploy/com.nuclear.cam-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nuclear.cam-server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/akietler/git/nuclear-eye/scripts/cam_server.py</string>
        <string>--port</string>
        <string>8085</string>
        <string>--device</string>
        <string>0</string>
        <string>--quality</string>
        <string>75</string>
        <string>--fps</string>
        <string>2</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/cam-server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/cam-server.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
