# deploy/nuclear-eye-vision-agent.service
[Unit]
Description=Nuclear Eye — Vision Agent (synthetic → real camera via SNAPSHOT_URL)
Documentation=https://github.com/marchange7/nuclear-eye
After=network.target nuclear-eye-alarm-grader.service
Wants=nuclear-eye-alarm-grader.service

[Service]
Type=simple
User=crew
Group=crew
WorkingDirectory=/home/crew/git/nuclear-eye
EnvironmentFile=-/home/crew/.config/nuclear-eye/env
Environment=RUST_LOG=info
Environment=HOUSE_SECURITY_CONFIG=/home/crew/git/nuclear-eye/config/security.toml
# Set SNAPSHOT_URL to enable real camera mode (e.g. http://m4:8085/snapshot)
# Set FASTVLM_URL to enable VLM captioning  (e.g. http://127.0.0.1:8091)
# Leave both unset to run synthetic test events
Environment=FASTVLM_URL=http://127.0.0.1:8091
ExecStart=/home/crew/git/nuclear-eye/target/release/vision_agent
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=60
StartLimitBurst=5
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
